# docker/cron.d/expensave-cron
0 * * * * www-data cd /opt/expensave && /usr/local/bin/php bin/console app:calendar:recalculate-balance > /proc/1/fd/1 2>/proc/1/fd/2
0 * * * * www-data cd /opt/expensave && /usr/local/bin/php bin/console gesdinet:jwt:clear > /proc/1/fd/1 2>/proc/1/fd/2
